# etc/supervisor.d/coredns.conf
[program:coredns]
command=/coredns -pidfile /var/run/coredns.pid -conf /etc/coredns/Corefile
