# etc/supervisor.d/docker-gen.conf
[program:docker-gen]
command=docker-gen -watch /etc/coredns/containers.tmpl /etc/coredns/containers
